# etc/supervisord.d/slurmstats.ini
[program:slurmstats]
command=/opt/qnib/slurm/bin/slurmstats.py -l --loop
autorestart=false
startsecs=15
startretries=10
redirect_stderr=true
stdout_logfile=syslog
